# systemd/ai-companion.service
[Unit]
Description=AI Companion Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/nyx/ai2d_chat
Environment=PATH=/home/nyx/ai2d_chat/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/nyx/ai2d_chat/.venv/bin/ai2d_chat server
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
